# site-cookbooks/uwsgi_app/templates/default/upstart_uwsgi.conf.erb
description "uwsgi app server"

start on runlevel [2345]
stop on runlevel [06]
 
exec /usr/local/bin/uwsgi --emperor /etc/uwsgi/vassals
start on startup
