# site-cookbooks/puma_app/templates/default/upstart_puma.conf.erb
description "Puma server for <%= @app_name %>"

start on runlevel [2345]
stop on runlevel [06]

setuid <%= @user %>
setgid <%= @user %>

respawn
respawn limit 3 30

env PATH=<%= "#{@shared_path}/ruby/bin:#{@shared_path}/bin:$PATH:/usr/local/bin:/usr/bin" %>
chdir <%= @current_path %>
exec bundle exec puma -C <%= @puma_config_path %>
